# deployment/systemd/salla-mcp.service
[Unit]
Description=Salla MCP Server
Documentation=https://github.com/jamalla/fastmcp-openapi-spec
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=salla-mcp
Group=salla-mcp
WorkingDirectory=/opt/salla-mcp
Environment="PATH=/opt/salla-mcp/venv/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONDONTWRITEBYTECODE=1"
EnvironmentFile=/etc/salla-mcp/.env
ExecStart=/opt/salla-mcp/venv/bin/python server.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=salla-mcp

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/salla-mcp

[Install]
WantedBy=multi-user.target
